# rootfs/etc/systemd/system/altitude-systemd-audit.timer
[Unit]
Description=Refresh Altitude systemd audit table

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h
AccuracySec=5min
Unit=altitude-systemd-audit.service

[Install]
WantedBy=timers.target
