# rootfs/etc/systemd/system/altitude-systemd-audit.timer
[Unit]
Description=Refresh Altitude systemd audit table

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=altitude-systemd-audit.service

[Install]
WantedBy=timers.target
