# rootfs/etc/systemd/system/altitude-systemd-audit.service
[Unit]
Description=Mirror systemd unit state into Altitude registry audit table
After=altitude-seed.service multi-user.target
Wants=altitude-seed.service

[Service]
Type=oneshot
ExecStart=/bin/systemd-audit refresh
TimeoutStartSec=30
Nice=10
IOSchedulingClass=idle
