# rootfs/etc/systemd/system/altitude-health.timer
[Unit]
Description=Refresh Altitude health snapshot

[Timer]
OnBootSec=45s
OnUnitActiveSec=5min
AccuracySec=30s
Unit=altitude-health.service

[Install]
WantedBy=timers.target
